# docker/supervisord.conf
; ============================================================
; Supervisord — manages Nginx, PHP-FPM, and Laravel queue worker
; All three processes run in a single container (Render free tier)
; ============================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
; Use /tmp so stale PID from a previous container never blocks startup
pidfile=/tmp/supervisord.pid

[unix_http_server]
; Use /tmp so stale socket from a previous container never blocks startup
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ------------------------------------------------------------
; PHP-FPM
; ------------------------------------------------------------
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/www.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/php-fpm.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stderr_logfile_maxbytes=5MB
; Give PHP-FPM time to start before Nginx tries to connect
startsecs=3

; ------------------------------------------------------------
; Nginx
; ------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/nginx-error.log
stderr_logfile_maxbytes=5MB
; Wait for PHP-FPM to be ready first
startsecs=5
startretries=3

; ------------------------------------------------------------
; Laravel queue worker (database driver)
; Uses --sleep to avoid hammering DB when queue is empty
; Uses --tries=3 and --timeout=90 (matches retry_after in config/queue.php)
; ------------------------------------------------------------
[program:laravel-worker]
command=php /var/www/html/artisan queue:work database --sleep=3 --tries=3 --timeout=90 --max-time=3600 --no-interaction
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=30
numprocs=1
stdout_logfile=/var/log/supervisor/worker.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/worker-error.log
stderr_logfile_maxbytes=5MB
; Start after Nginx is up
startsecs=8
stopwaitsecs=120
; Restart the worker when it exits (e.g. after --max-time)
exitcodes=0,1
